# 1.Base/default/debian-core/etc/protect/template.core-keys/config/rc.local
#!/bin/sh

/sbin/fstrim --quiet /rw

/sbin/shutdown +2

exit 0
