# 2.Network/default.first/debian-core/etc/protect/template.core-tor/config/rc.local
#!/bin/sh

RETURN="
"
NS1="10.139.1.1"
NS2="10.139.1.2"
if [ -e /var/run/qubes/qubes-ns ] ; then
    . /var/run/qubes/qubes-ns
fi
echo "nameserver ${NS1}${RETURN}nameserver ${NS2}" > /run/resolv.conf

/sbin/fstrim --quiet /rw

/bin/systemctl enable qubes-iptables
/bin/systemctl start qubes-iptables

. /rw/config/qubes-firewall-user-script

/bin/systemctl enable liteqube-tor-monitor
/bin/systemctl start liteqube-tor-monitor

/usr/bin/qrexec-client-vm core-net liteqube.WifiRequestAP
