# domU/etc/protect/template.core-xorg/config/rc.local
#!/bin/sh

/bin/systemctl start qubes-gui-agent
