# 7.Audio/custom/debian-core/etc/protect/template.core-sound/config/rc.local
#!/bin/sh

systemctl start liteqube-pulseaudio-hardware
systemctl start liteqube-pulseaudio-daemon
systemctl start liteqube-pulseaudio-monitor

sleep 3
[ -x "/rw/config/liteqube-settings" ] && . /rw/config/liteqube-settings
[ -n "${INITIAL_OUTPUT_UNMUTE}" ] && sudo -u user /etc/qubes-rpc/liteqube.SoundVolume UNMUTE || sudo -u user /etc/qubes-rpc/liteqube.SoundVolume MUTE
[ -n "${INITIAL_OUTPUT_VOLUME}" ] && sudo -u user /etc/qubes-rpc/liteqube.SoundVolume "${INITIAL_OUTPUT_VOLUME}"
