# domU/etc/protect/template.core-mail/config/rc.local
#!/bin/sh

/bin/systemctl enable getmail.timer
/bin/systemctl start getmail.timer
