# domU/etc/protect/template.core-tor/config/rc.local
#!/bin/sh

/bin/systemctl enable tor-monitor
/bin/systemctl start tor-monitor
