# 2.Network/default.fw-linux/debian-core/etc/protect/template.fw-tor/config/rc.local
#!/bin/sh

/bin/systemctl enable qubes-iptables
/bin/systemctl start qubes-iptables

/bin/systemctl enable qubes-network
/bin/systemctl start qubes-network
