# 1.Base/default/debian-core/etc/protect/template.core-dvm/config/rc.local
#!/bin/sh

chattr -R -i /rw/QUARANTINE
chattr -R -i /rw/home
chattr -R -i /rw/config

rm -rf /rw/QUARANTINE
rm -rf /rw/home/*
rm -rf /rw/config/*

shutdown now
