# Contrib.SmartCard/default/debian-core/lib/systemd/system/pcscd.service.d/30_qubes.conf
[Unit]
ConditionPathExists=/var/run/qubes-service/qubes-pkcs11
After=qubes-sysinit.service
