# 7.Audio/default/debian-core/etc/protect/template.core-sound/config/rc.local
#!/bin/sh

systemctl start liteqube-pulseaudio-hardware
systemctl start liteqube-pulseaudio-daemon
systemctl start liteqube-pulseaudio-monitor
