# domU/lib/systemd/system/getmail.timer
[Unit]
Description=Run getmail every 15 minutes

[Timer]
OnStartupSec=1sec
OnUnitInactiveSec=15min

[Install]
WantedBy=timers.target