# 1.Base/default/debian-core/etc/rc.local
#!/bin/sh -e

if [ -e /rw/QUARANTINE ] ; then
    sleep 3
    /usr/bin/qrexec-client-vm dom0 "liteqube.Error+Files quarantined during boot"
fi

exit 0
