# 1.Base/default/debian-core/etc/rc.local
#!/bin/sh -e

NS1="10.139.1.1"
NS2="10.139.1.2"
if [ -e /var/run/qubes/qubes-ns ] ; then
    . /var/run/qubes/qubes-ns
fi
echo -e "nameserver ${NS1}\nnameserver ${NS2}" > /run/resolv.conf

/sbin/fstrim -a

if [ -e /rw/QUARANTINE ] ; then
    sleep 3
    /usr/bin/qrexec-client-vm dom0 "liteqube.Error+Files quarantined during boot"
fi

exit 0
