# 1.Base/default/debian-core/etc/protect/template.core-xorg/config/rc.local
#!/bin/sh

/bin/systemctl start qubes-gui-agent

/bin/systemctl enable liteqube-shutdown-inactive
/bin/systemctl start liteqube-shutdown-inactive
