# 3.USB/default/debian-core/etc/protect/template.core-usb/config/rc.local
#!/bin/sh

/bin/systemctl enable usbguard
/bin/systemctl start usbguard
