# 5.Storage/default/debian-core/etc/protect/template.core-decrypt/config/rc.local
#!/bin/sh

/bin/systemctl enable dm-event.socket
/bin/systemctl start dm-event.socket
